# systemd/captivity.service
[Unit]
Description=Captivity — Autonomous Captive Portal Login Client
Documentation=https://github.com/gaminization/captivity
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
# Use the pip-installed CLI entry point
# For system-wide install, change to: /usr/local/bin/captivity daemon
ExecStart=%h/.local/bin/captivity daemon
Restart=on-failure
RestartSec=10
# Allow reading credentials from user keyring
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
# Load network config (optional)
EnvironmentFile=-%h/.config/captivity/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=captivity

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/tmp %h/.local/share/captivity

[Install]
WantedBy=default.target
