# systemd/captivity.service
[Unit]
Description=Captivity — Autonomous Captive Portal Login Client
After=NetworkManager.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=garvarora
Group=garvarora
WorkingDirectory=/home/garvarora

ExecStart=/home/garvarora/.local/bin/captivity daemon

Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/garvarora/.Xauthority

# Restart quickly on crash; the startup-reconciliation window in the daemon
# makes each restart self-correcting even if the network is not ready yet.
Restart=always
RestartSec=3

StartLimitBurst=10
StartLimitIntervalSec=60

[Install]
WantedBy=multi-user.target
